# rootfs/etc/s6-overlay/s6-rc.d/socks-to-http-proxy/run
#!/command/with-contenv sh
exec ${BASE_PATH}/socks-to-http-proxy/sthp-linux --port 2335 --listen-ip 0.0.0.0 --socks-address 127.0.0.1:2334 --http-basic ${PROXY_USER}:${PROXY_PASS}